# digitalocean/config/ztf-kafka.target
[Unit]
Description=ZTF Kafka Alert Broker Cluster
AllowIsolate=no
